# memavaild.service.in
[Unit]
Description=Improve responsiveness during heavy swapping: keep amount of available memory
Documentation=https://github.com/hakavlad/memavaild

[Service]
ExecStart=:TARGET_SBINDIR:/memavaild -c :TARGET_SYSCONFDIR:/memavaild.conf
Restart=on-failure
User=memavaild
Slice=critical.slice
AmbientCapabilities=CAP_IPC_LOCK CAP_DAC_OVERRIDE
TasksMax=1
MemoryMax=25M
MemorySwapMax=1
PrivateTmp=true
ProtectHome=true
ProtectClock=true
ProtectSystem=strict
ProtectHostname=true
ProtectKernelModules=true
ProtectKernelTunables=true
MemoryDenyWriteExecute=yes
RestrictNamespaces=yes
RestrictRealtime=yes
PrivateDevices=true
PrivateNetwork=true
LockPersonality=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
