# memavaild.service.in
[Unit]
Description=Keep amount of available memory
Documentation=https://github.com/hakavlad/memavaild

[Service]
ExecStart=:TARGET_SBINDIR:/memavaild
Restart=always
RestartSec=0
User=memavaild
AmbientCapabilities=CAP_IPC_LOCK CAP_DAC_OVERRIDE
TasksMax=1
MemoryMax=25M
MemorySwapMax=1M
ProtectHome=true
ProtectSystem=strict
MemoryDenyWriteExecute=yes
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictNamespaces=yes
RestrictRealtime=yes
PrivateDevices=true
LockPersonality=yes
PrivateNetwork=true
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
