# memavaild.service.in
[Unit]
Description=memavaild

[Service]
ExecStart=:TARGET_SBINDIR:/memavaild
Restart=always
RestartSec=0
RestrictRealtime=yes
TasksMax=1
MemoryMax=20M
MemorySwapMax=1M
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
PrivateDevices=true
MemoryDenyWriteExecute=yes
RestrictNamespaces=yes
LockPersonality=yes
PrivateNetwork=true

[Install]
WantedBy=multi-user.target
